# fprintd.service.d/r503-arduino-legacy.conf
[Service]
Environment="R503_ARDUINO_DEVICE=1"
